# ansible/roles/stack/files/etc/cron.d/track_coin.conf
01 * * * * root sh /srv/app.sh
*/1 22 * * * root sh /srv/app.sh
*/1 23 * * * root sh /srv/app.sh
